# src/deploy/NVA_build/supervisord.orig
#!/bin/sh

# chkconfig: 345 15 15

. /etc/rc.d/init.d/functions

SUPERVISORD="/usr/bin/supervisord_orig"
PIDFILE="/tmp/supervisord.pid"

start() {
    if [ ! -x "$SUPERVISORD" ]; then
        echo "$SUPERVISORD is not executable."
        logger -p local0.info "$SUPERVISORD is not executable."
        exit 1
    fi
    echo "Starting ..."
    logger -p local0.info "Starting ..."
    ulimit -n 102400; $SUPERVISORD --pidfile $PIDFILE

    return $?
}

stop() {
    echo "Stopping ..."
    logger -p local0.info "Stopping ..."
    local supervisord_orig_pid=$(pgrep -lf $SUPERVISORD | awk '{print $1}')
    echo "supervisord_orig_pid=$supervisord_orig_pid"
    logger -p local0.info "supervisord_orig_pid=$supervisord_orig_pid"
    kill -9 ${supervisord_orig_pid}

    local web_server_pid=$(pgrep -lf web_server.js | awk '{print $1}')
    echo "web_server_pid=$web_server_pid"
    logger -p local0.info "web_server_pid=$web_server_pid"
    kill -9 ${web_server_pid}

    local bg_workers_pid=$(pgrep -lf bg_workers.js | awk '{print $1}')
    echo "bg_workers_pid=$bg_workers_pid"
    logger -p local0.info "bg_workers_pid=$bg_workers_pid"
    kill -9 ${bg_workers_pid}

    local hosted_agents_starter_pid=$(pgrep -lf hosted_agents_starter.js | awk '{print $1}')
    echo "hosted_agents_starter_pid=$hosted_agents_starter_pid"
    logger -p local0.info "hosted_agents_starter_pid=$hosted_agents_starter_pid"
    kill -9 ${hosted_agents_starter_pid}

    local s3rver_starter_pids=$(pgrep -lf s3rver_starter.js | awk '{print $1}')
    echo "s3rver_starter_pids=$s3rver_starter_pids"
    logger -p local0.info "s3rver_starter_pids=$s3rver_starter_pids"
    kill -9 ${s3rver_starter_pids}

    local mongo_pids=$(pgrep -lf mongo | awk '{print $1}')
    echo "mongo_pids=$mongo_pids"
    logger -p local0.info "mongo_pids=$mongo_pids"
    kill -9 ${mongo_pids}

    local p=$(cat $PIDFILE)
    kill -9 ${p}
    [ $? -eq 0 ] && rm -f $PIDFILE
    return $retval
}

case $1 in
    start)
        start
    ;;
    stop)
        stop
    ;;
    restart)
        stop
        sleep 1
        start
    ;;
    *)
        echo "$0 start|stop|restart"
        exit 2
    ;;
esac
